# .github/actions/deploy/cloud-init/write_files/etc/systemd/system/mapineda48-loki.service
# /etc/systemd/system/loki.service
[Unit]
Description=Grafana Loki
After=network.target

[Service]
ExecStart=/usr/local/bin/loki --config.file=/etc/loki/config.yml
Restart=always

[Install]
WantedBy=multi-user.target
